# infra/systemd/kiln-firecracker-host-manager.service
[Unit]
Description=Kiln Firecracker Host Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=__KILN_APP_DIR__/packages/runner
EnvironmentFile=/etc/kiln/app.env
ExecStart=/usr/bin/env node dist/sandbox/host-manager.js
Restart=always
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
